# recipes-forwarders/restconfs/files/restconfs.socket
[Unit]
Description=Socket for restconfs

[Socket]
ListenStream=443

[Install]
WantedBy=sockets.target
