# recipes-forwarders/usbconfs/files/confs-usbstorage@.path
[Unit]
Description=confs mount point monitoring
#BindsTo=%i.mount

[Path]
PathExists=/%I
Unit=confs-usbstorage@%i.service
